# roles/dgimonitorrsync/templates/dgichangemonitor.service.j2
[Unit]
Description= Starts Data Grid Rsync
Documentation= man:systemd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "echo 'rsync' >> /tmp/arlindo.txt"
