# roles/dgimonitorrsync/templates/dgichangemonitor.path.j2
[Unit]
Description=Triggers the service that monitors data grid file changes.
Documentation=man:systemd.path

[Path]
PathModified=/tmp/testchanges/

[Install]
WantedBy=multi-user.target
